# AWDLControl/AWDLControlHelper/com.awdlcontrol.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.awdlcontrol.helper</string>

    <key>MachServices</key>
    <dict>
        <key>com.awdlcontrol.helper</key>
        <true/>
    </dict>

    <!-- Path to the helper executable -->
    <key>Program</key>
    <string>Contents/Library/LaunchServices/com.awdlcontrol.helper</string>

    <!-- SMAppService configuration -->
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.awdlcontrol.app</string>
    </array>
</dict>
</plist>
